# bot.conf
[program:telegrambot]
directory=/home/ubuntu/TelegramBots/Footballbot
command=python3 bot.py
environment=HOME="/home/ubuntu",USER="ubuntu"

autostart=true
autorestart=true